# deploy/sandbox/yanban-sandbox-broker.service.example
[Unit]
Description=Yanban governed Sandbox Broker
After=network-online.target mysql.service

[Service]
Type=simple
User=yanban-sandbox
Group=yanban-sandbox
SupplementaryGroups=kvm
EnvironmentFile=/etc/yanban/sandbox-broker.env
ExecStart=/usr/bin/java -jar /opt/yanban/yanban-sandbox-broker.jar
WorkingDirectory=/var/lib/yanban-sandbox
StateDirectory=yanban-sandbox
StateDirectoryMode=0700
Environment=HOME=/var/lib/yanban-sandbox
Environment=XDG_CONFIG_HOME=/var/lib/yanban-sandbox/config
Environment=XDG_DATA_HOME=/var/lib/yanban-sandbox/data
Environment=XDG_STATE_HOME=/var/lib/yanban-sandbox/state
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=no
DevicePolicy=closed
DeviceAllow=/dev/kvm rw
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/yanban-sandbox
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0077
Restart=on-failure

[Install]
WantedBy=multi-user.target
